# root/etc/s6-overlay/s6-rc.d/svc-exporter/run
#!/usr/bin/with-contenv bash
exec python3 /exporter/cli.py run