# root/etc/s6-overlay/s6-rc.d/svc-exporter/run
#!/command/with-contenv bash
python3 /exporter/cli.py run